# misc/systemd/nix-daemon.service.in
[Unit]
Description=Nix Daemon

[Service]
User=nix-daemon
EnvironmentFile=/etc/sysconfig/nix-daemon
ExecStart=@@bindir@/nix-daemon nix-daemon --daemon
KillMode=process
