# misc/systemd/nix-daemon.service.in
[Unit]
Description=Nix Daemon

[Service]
User=nix-daemon
RuntimeDirectory=@runstatedir@
EnvironmentFile=/etc/sysconfig/nix-daemon
ExecStart=@@bindir@/nix-daemon nix-daemon --daemon
KillMode=process

[Install]
WantedBy=multi-user.target 
