# misc/systemd/nix-daemon.service.in
[Unit]
Description=Nix Daemon
RequiresMountsFor=@storedir@
RequiresMountsFor=@localstatedir@
ConditionPathIsReadWrite=@localstatedir@/nix/daemon-socket

[Service]
User=nix-daemon
ExecStart=@@bindir@/nix-daemon nix-daemon --daemon
KillMode=process
